# supervisord.conf
[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
user=root

[program:gophernicus]
command=/gophernicus.sh

[program:gostatic]
command=/goStatic